# scripts/linux/nv-backend.service
[Unit]
Description=NV Local Voice Lead Closer backend
After=network-online.target lifeos-ollama.service hermes-gateway.service
Wants=lifeos-ollama.service hermes-gateway.service

[Service]
Type=simple
WorkingDirectory=/home/dell/ram/nv/backend
EnvironmentFile=/home/dell/ram/nv/.env
ExecStart=/home/dell/ram/nv/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8090
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
